# wifi-check.cron
* * * * * /home/pi/wifi-check.sh >> /home/pi/.log/network.log 2>&1
